# s6-overlay/s6-rc.d/dote/run
#!/command/with-contenv bash
exec /command/s6-setuidgid dote /usr/local/bin/dote $DOTE_OPTS 2>&1
